# deployments/com.winfertility.scout.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.winfertility.scout</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/expo/fertility-sense/.venv/bin/python</string>
        <string>-m</string>
        <string>fertility_sense</string>
        <string>scout</string>
        <string>--once</string>
    </array>
    <key>WorkingDirectory</key><string>/Users/expo/fertility-sense</string>
    <key>StartInterval</key><integer>21600</integer>
    <key>StandardOutPath</key><string>/Users/expo/fertility-sense/data/logs/scout.log</string>
    <key>StandardErrorPath</key><string>/Users/expo/fertility-sense/data/logs/scout.err</string>
</dict>
</plist>
