# deployments/com.winfertility.sequence-run.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key><string>com.winfertility.sequence-run</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/expo/fertility-sense/.venv/bin/python</string>
        <string>-m</string>
        <string>fertility_sense</string>
        <string>sequence</string>
        <string>run</string>
    </array>
    <key>WorkingDirectory</key><string>/Users/expo/fertility-sense</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>9</integer>
        <key>Minute</key><integer>0</integer>
    </dict>
    <key>StandardOutPath</key><string>/Users/expo/fertility-sense/data/logs/sequence-run.log</string>
    <key>StandardErrorPath</key><string>/Users/expo/fertility-sense/data/logs/sequence-run.err</string>
</dict>
</plist>
